# build/fk-w3-tg-aw-reporting/deploy/cron.d/fk-w3-tg-aw-reporting
MAILTO='amol.sharma@flipkart.com,santhana.kumar@flipkart.com,vignesh.raghuraman@flipkart.com'

# Metrics for all the accounts
#00 13 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchMetricDelayedReports.sh
30 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchMetricReports.sh

# Attribute data download by account size
40 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeMisc.sh
50 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeMobile.sh
00 19 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeApparel.sh
10 19 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeShoes.sh
20 19 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeToys.sh
30 19 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeBooks.sh
40 19 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeDisplay.sh
