# build/fk-w3-tg-aw-reporting/deploy/cron.d/fk-w3-tg-aw-reporting
MAILTO='tg-sem@flipkart.com'

0 0 2 * * /usr/share/fk-w3-tg-aw-reporting/bin/fetchReports.sh

