# build/fk-w3-tg-aw-reporting/deploy/cron.d/fk-w3-tg-aw-reporting
MAILTO='anurag.laddha@flipkart.com'

#15 15 * * 6 root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchMetricReports.sh
15 17 * * 6 root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeReports.sh


