# build/fk-w3-tg-aw-reporting/deploy/cron.d/fk-w3-tg-aw-reporting
MAILTO='alok.reddy@flipkart.com,vignesh.raghuraman@flipkart.com'

# Metrics for all the accounts
30 16 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchMetricDelayedReports.sh
45 16 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchMetricReports.sh

# Attribute data download by account size
00 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeMisc.sh
10 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeMobile.sh
20 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeApparel.sh
30 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeShoes.sh
40 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeToys.sh
50 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeBooks.sh
00 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeDisplay.sh
10 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeHomeAndKitchen.sh
