# build/fk-w3-tg-aw-reporting/deploy/cron.d/fk-w3-tg-aw-reporting
MAILTO='alok.reddy@flipkart.com,vignesh.raghuraman@flipkart.com'

# Metrics for all the accounts
15 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchMetricDelayedReports.sh
30 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchMetricReports.sh

# Attribute data download by account size
45 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeMisc.sh
55 17 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeMobile.sh
05 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeApparel.sh
15 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeShoes.sh
25 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeToys.sh
35 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeBooks.sh
45 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeDisplay.sh
55 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeHomeAndKitchen.sh
