# build/fk-w3-tg-aw-reporting/deploy/cron.d/fk-w3-tg-aw-reporting
MAILTO='anurag.laddha@flipkart.com'

#10 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchMetricReports.sh
50 1 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeReports.sh


