# build/fk-w3-tg-aw-reporting/deploy/cron.d/fk-w3-tg-aw-reporting
MAILTO='santhana.kumar@flipkart.com'

0 18 * * * root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchReports.sh

