# build/fk-w3-tg-aw-reporting/deploy/cron.d/fk-w3-tg-aw-reporting
MAILTO='anurag.laddha@flipkart.com'

#metrics for all the accounts
#15 15 * * 6 root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchMetricReports.sh

#attribute data download by account size
40 11 * * 7 root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeMisc.sh
50 11 * * 7 root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeMobile.sh
59 11 * * 7 root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeApparel.sh
10 12 * * 7 root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeShoes.sh
20 12 * * 7 root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeToys.sh
30 12 * * 7 root sh /usr/share/fk-w3-tg-aw-reporting/bin/fetchAttributeBooks.sh


